# tradedeck/deploy/tradedeck.service
[Unit]
Description=TradeDeck v2 Docker Stack
Requires=docker.service
After=docker.service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/tradedeck-v2
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
User=root
EnvironmentFile=/home/ubuntu/tradedeck-v2/.env

[Install]
WantedBy=multi-user.target
